# deploy/gebrauchtwaffen.service
# Systemd service file for Gilbert's Yoga Helper
#
# Installation:
#   sudo cp deploy/yoga-helper.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable yoga-helper
#   sudo systemctl start yoga-helper
#
# Management:
#   sudo systemctl status yoga-helper
#   sudo systemctl stop yoga-helper
#   sudo systemctl restart yoga-helper
#   sudo journalctl -u yoga-helper -f
#
# Configuration:
#   Adjust User, WorkingDirectory, and ExecStart paths as needed

[Unit]
Description=Gilbert's Yoga Helper - Web Aggregator
Documentation=https://github.com/your-repo/yoga-helper-aggregator
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/yoga-helper_aggregator

# Application startup
ExecStart=/home/pi/yoga-helper_aggregator/.venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8000

# Restart configuration
Restart=always
RestartSec=10

# Stop gracefully
TimeoutStopSec=30
KillMode=mixed

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=LOG_LEVEL=INFO
Environment=USE_JOURNALD=true

# Security hardening (optional, adjust as needed)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/yoga-helper_aggregator/data
ReadWritePaths=/home/pi/yoga-helper_aggregator/logs

# Resource limits (adjust for Pi)
MemoryMax=512M
CPUQuota=80%

[Install]
WantedBy=multi-user.target
